# opencode-telegram.service
[Unit]
Description=OpenCode Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/opencode-telegram
ExecStart=/usr/bin/python3 /opt/opencode-telegram/bot.py
Restart=always
RestartSec=10

# Environment variables - edit these or use an EnvironmentFile
Environment=TELEGRAM_BOT_TOKEN=your-telegram-bot-token-here
Environment=TELEGRAM_ALLOWED_USER_IDS=your-telegram-user-id
Environment=OPENCODE_PATH=/root/.opencode/bin/opencode
Environment=OPENCODE_DEFAULT_CWD=/root
Environment=OPENCODE_LOG_FILE=/var/log/opencode-telegram/bot.log
Environment=OPENCODE_SESSIONS_FILE=/var/lib/opencode-telegram/sessions.json

# Alternatively, use an environment file:
# EnvironmentFile=/opt/opencode-telegram/.env

[Install]
WantedBy=multi-user.target
