# opencode-serve.service
[Unit]
Description=OpenCode Server Daemon
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/root/.opencode/bin/opencode serve --port 8080 --hostname 127.0.0.1
Restart=always
RestartSec=5

# Environment for OpenCode
Environment=HOME=/root

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=opencode-serve

[Install]
WantedBy=multi-user.target
